# etc/cron.d/app
MAILTO=egor.orlov@ostrovok.ru

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user command

# daily
0 * * * *  root . /etc/container_environment.sh && cd /opt/app/ && run-one ./manage.py hourly_check --traceback
